# systemd/wifi-manager.service
[Unit]
Description=WiFi Auto Reconnect Script
After=network.target

[Service]
ExecStart=/opt/wifi-manager/auto_connect.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
